# runit/dwmstatus/run
#!/bin/sh
export DISPLAY=":0"
exec 2>&1
exec dwmstatus
