# runit/foo/run
#!/bin/sh -e
exec 2>&1
exec chpst -u risto /home/risto/scripts/runit/foo/foo.sh
