# runit/emacsd/run
#!/bin/sh

export HOME=/home/risto
cd $HOME

exec 2>&1
exec emacs --fg-daemon=emacsd
